# services/sentinel-wifi.service
[Unit]
Description=SENTINEL WiFi Probe Capture Engine
Documentation=file:///opt/sentinel/README.md
After=network.target
Wants=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m capture.wifi_capture
WorkingDirectory=/opt/sentinel
Restart=on-failure
RestartSec=10

# Run as root for raw socket access
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-wifi

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/sentinel/data /opt/sentinel/logs
ProtectHome=yes
NoNewPrivileges=no
PrivateTmp=yes

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentinel

[Install]
WantedBy=multi-user.target
